# package-main/usr/lib/systemd/system/kfocus-rollback-timed-snapshot.service
[Unit]
Description=FocusRx Rollback Timed Service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/lib/kfocus/bin/kfocus-focusrx-system -s
