# package-nvidia/debian/kfocus-nvidia.prerm
#!/bin/bash

set -e;

deb-systemd-helper disable nvidia-suspend.service;
deb-systemd-helper disable nvidia-hibernate.service;
deb-systemd-helper disable nvidia-resume.service;

# Tag to allow some debhelper commands to inject relevant code
#DEBHELPER#
